# data/com.github.darkoverlordofdata.bada-bing.cron
0 */6 * * * com.github.darkoverlordofdata.bada-bing --update >/dev/null 2>&1